# flask/update.cron
# cron job that updates SBOLExplorer at the zero of every hour
0 * * * * curl -X GET "localhost:13162/update" >> /tmp/update.cron.output
